# setup/supervisord.conf
[supervisord]
nodaemon=true

[program:install]
command=/src/zodiac/setup/install.sh
directory=/src/zodiac/
user=zodiac
