# setup/supervisord.conf
[supervisord]
nodaemon=true
user=zodiac

[program:install]
command=/src/zodiac/setup/install.sh
directory=/src/zodiac/
